# disk_space/crontab
0 1 * * * timeout 6h /homes/muffato/workspace/mysql-monitor/disk_space/cron.sh
